# misc/osx/LaunchAgents/kirbies.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>kirbies.tunnel</string>
	<key>ProgramArguments</key>
	<array>
		<string>ssh</string>
		<string>-NTC</string>
		<string>-o</string>
		<string>ExitOnForwardFailure=yes</string>
		<string>-R</string>
		<string>10022:localhost:22</string>
		<string>kirbies</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
</dict>
</plist>
